# rclone-systemd-ivey/rclone-backup@.timer
[Unit]
Description=Backup files timer

[Timer]
OnActiveSec=300
OnCalendar=*-*-* 2/6:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target

