# cloudflare-refresh/cloudflare-refresh.timer
[Unit]
Description=Refresh cloudflare timer

[Timer]
OnActiveSec=300
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target