# lan-firewall/lan-firewall.timer
[Unit]
Description=Refresh lan firewall config timer

[Timer]
OnActiveSec=0
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target